# templates/default/jetbrains.teamcity.BuildAgent.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Debug</key>
	<false/>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>jetbrains.teamcity.BuildAgent</string>
	<key>LowPriorityBackgroundIO</key>
	<false/>
	<key>LowPriorityIO</key>
	<false/>
	<key>OnDemand</key>
	<false/>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>--login</string>
		<string>-c</string>
		<string><%= @launch_command %></string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>logs/launchd.err.log</string>
	<key>StandardOutPath</key>
	<string><%= @stdout_path %></string>
	<key>WorkingDirectory</key>
	<string><%= @working_dir %></string>
</dict>
</plist>
